# obsidian-blog.service
[Unit]
Description=Obsidian Blog Generator
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html
ExecStart=/usr/bin/node /var/www/html/index.js watch
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=obsidian-blog

# 환경 변수
Environment=NODE_ENV=production

# 보안 설정
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/html

# 리소스 제한
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target 